# server/deploy/runordie.service
# systemd-Unit für das RUN OR DIE Backend — Relay + Konten-API (Frontend: GitHub Pages)
#   sudo systemctl daemon-reload && sudo systemctl enable --now runordie
[Unit]
Description=RUN OR DIE backend (WebSocket relay + accounts API)
After=network.target

[Service]
User=fab
WorkingDirectory=/opt/runordie/server
ExecStart=/usr/bin/env node server.js
Environment=PORT=18931
# Frontend-Origin (GitHub Pages): erlaubt für WebSocket UND CORS. localhost ist
# für Entwicklung ohnehin immer erlaubt. Bei eigener Domain hier ergänzen.
Environment=ORIGINS=https://fl4p.github.io
Environment=CORS_ORIGINS=https://fl4p.github.io
# Hinter nginx: echte Client-IP aus X-Forwarded-For lesen (sonst zählen alle
# Rate-Limits gegen 127.0.0.1). nginx MUSS den Header setzen (tut es im Snippet).
Environment=TRUST_PROXY=1
# Registrierungs-Limit deaktiviert (0 = aus); zum Wiedereinschalten z.B. 20 setzen.
Environment=REG_MAX=0
Restart=always
RestartSec=3
# Härtung: das Relay braucht weder Schreibrechte noch neue Privilegien
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
